# daemon/infrafid.service
[Unit]
Description=InfraFi IR WiFi Credential Receiver
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/infrafid --foreground
Restart=always
RestartSec=5
User=root
Group=root

# Logging goes to syslog/journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=infrafid

[Install]
WantedBy=multi-user.target
