# daemon/wifird.service
[Unit]
Description=Wi-Fir IR WiFi Credential Receiver
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/wifird --foreground
Restart=always
RestartSec=5
User=root
Group=root

# Logging goes to syslog/journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wifird

[Install]
WantedBy=multi-user.target
